# tags/1.3.99.1/demos/XrdMon/rpm/configs/etc/cron.d/gled-xrdmon-check
SHELL=/bin/sh
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root

LOGFILE=/var/log/gled/xrdmon-check.log

*/5 * * * * root /sbin/service gled-xrdmon _status_running >/dev/null 2>&1; [ $? -ne 0 ] && (echo "[`date`] gled-xrdmon-check: some collectors supposed to be running are stopped (crashed?), will be restarted"  && /sbin/service gled-xrdmon _recovery_restart) | tee -a $LOGFILE | mail -s "Gled XrdMon collectors checker" $MAILTO
